# deploy/templates/app.service.j2
[Unit]
Description={{ app_service_name }} API service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
EnvironmentFile={{ app_dir }}/.env
ExecStart={{ app_dir }}/dist/server
Restart=on-failure
RestartSec=3
TimeoutStartSec=60
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target